# images/vm/assets/unisonos-firstboot.service
[Unit]
Description=UnisonOS first-boot provisioning
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/unisonos-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target

