# systemd/unison-consent.service
[Unit]
Description=Unison Consent Service
Documentation=https://docs.unison.ai
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=unison
Group=unison
WorkingDirectory=/opt/unison/unison-consent

# Environment
EnvironmentFile=/opt/unison/.env
Environment="PYTHONPATH=/opt/unison/unison-consent/src"
Environment="UNISON_CONSENT_PORT=7072"
Environment="UNISON_CONSENT_HOST=0.0.0.0"
Environment="UNISON_CONSENT_KEYS_DIR=/var/lib/unison/keys"

# Start command
ExecStart=/usr/bin/python3.12 /opt/unison/unison-consent/src/main.py

# Restart policy
Restart=on-failure
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=100%

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/unison /var/log/unison

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unison-consent

[Install]
WantedBy=multi-user.target
